# puppet/modules/tomcat/templates/tomcat.sh.erb
#!/bin/sh
#
# mulgara:	Mulgara Daemon
#
# chkconfig: 35 99 1
#
# description: Starts and stops the Tomcat service
#
# processname: mulgara
# pidfile: /var/run/tomcat.pid
#
REPO_DIR=/opt/repository

if [ -x $REPO_DIR/cdr-env.sh ]; then
	. $REPO_DIR/cdr-env.sh
else
	echo "CDR environment script missing from REPO_DIR: ${REPO_DIR}/cdr-env.sh"
	exit 0;
fi


# Source function library.
. /etc/rc.d/init.d/functions

start() {
# 	next line removed by Greg as redundant
#	export JAVA_OPTS="-Xms128m -Xmx1024m"
	echo "Starting CDR Tomcat Server: "
        runuser <%= @user %> -p -c "$CATALINA_HOME/bin/startup.sh"
        ret=$?

        if [ $ret -eq 0 ]; then
		echo_success
        else
        	echo_failure
        fi
}

stop() {
	echo -n "Stopping CDR Tomcat Server: "
        runuser <%= @user %> -p -c "$CATALINA_HOME/bin/shutdown.sh"
	ret=$?
	if [ $ret -eq 0 ]; then
		echo_success
        else
		echo_failure
        fi
}

### main logic ###
case "$1" in
  start)
        start
        ;;
  stop)
        stop
        ;;
  status)
        status FOO
        ;;
  restart|reload|condrestart)
        stop
	sleep 10
        start
        ;;
  *)
        echo $"Usage: $0 {start|stop|restart|reload|status}"
        exit 1
esac

exit 0
